# tools/smartfield.service
[Unit]
Description=SmartField Rover Backend
After=network.target

[Service]
# User running the service (Change if not 'pratik')
User=pratik
Group=pratik

# Working Directory (Change if you cloned elsewhere)
# Assumes: /home/pratik/smartfield-companion
WorkingDirectory=/home/pratik/smartfield-companion

# Environment variables (Optional)
# ENV=production

# Command to start the server
# Using python directly to run the uvicorn call in app.py or uvicorn directly
# Adjust path to python if using a venv: /home/pratik/smartfield-companion/venv/bin/python
ExecStart=/usr/bin/python3 -m uvicorn backend.main:app --host 0.0.0.0 --port 8000

# Restart on crash
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
